# zookeeper/kafka/meta/filelist.json
{
   "kafka-server" : {
    "name": "kafka-server",
    "version": "1.0.1+kafka3.1.0",
    "files" : {
      "etc/rc.d/init.d/kafka-server" : {}
    }
  }
  , "kafka-mirror-maker" : {
    "name": "kafka-mirror-maker",
    "version": "1.0.1+kafka3.1.0",
    "files" : {
      "etc/rc.d/init.d/kafka-mirror-maker" : {}
    }
  }
  , "kafka" : {
    "name": "kafka",
    "version": "1.0.1+kafka3.1.0",
    "files" : {
      "etc/default/kafka" : {},
      "etc/kafka/conf.dist/connect-console-sink.properties" : {},
      "etc/kafka/conf.dist/connect-console-source.properties" : {},
      "etc/kafka/conf.dist/connect-distributed.properties" : {},
      "etc/kafka/conf.dist/connect-file-sink.properties" : {},
      "etc/kafka/conf.dist/connect-file-source.properties" : {},
      "etc/kafka/conf.dist/connect-log4j.properties" : {},
      "etc/kafka/conf.dist/connect-standalone.properties" : {},
      "etc/kafka/conf.dist/log4j.properties" : {},
      "etc/kafka/conf.dist/server.properties" : {},
      "etc/kafka/conf.dist/tools-log4j.properties" : {},
      "bin/kafka-broker-api-versions" : {},
      "bin/kafka-config" : {},
      "bin/kafka-console-consumer" : {},
      "bin/kafka-console-producer" : {},
      "bin/kafka-consumer-groups" : {},
      "bin/kafka-consumer-offset-checker" : {},
      "bin/kafka-consumer-perf-test" : {},
      "bin/kafka-preferred-replica-election" : {},
      "bin/kafka-producer-perf-test" : {},
      "bin/kafka-reassign-partitions" : {},
      "bin/kafka-run-class" : {},
      "bin/kafka-sentry" : {},
      "bin/kafka-topics" : {},
      "lib/kafka" : {},
      "lib/kafka/LICENSE" : {},
      "lib/kafka/NOTICE" : {},
      "lib/kafka/bin" : {},
      "lib/kafka/bin/connect-distributed.sh" : {},
      "lib/kafka/bin/connect-standalone.sh" : {},
      "lib/kafka/bin/kafka-acls.sh" : {},
      "lib/kafka/bin/kafka-broker-api-versions.sh" : {},
      "lib/kafka/bin/kafka-configs.sh" : {},
      "lib/kafka/bin/kafka-console-consumer.sh" : {},
      "lib/kafka/bin/kafka-console-producer.sh" : {},
      "lib/kafka/bin/kafka-consumer-groups.sh" : {},
      "lib/kafka/bin/kafka-consumer-offset-checker.sh" : {},
      "lib/kafka/bin/kafka-consumer-perf-test.sh" : {},
      "lib/kafka/bin/kafka-delete-records.sh" : {},
      "lib/kafka/bin/kafka-log-dirs.sh" : {},
      "lib/kafka/bin/kafka-mirror-maker.sh" : {},
      "lib/kafka/bin/kafka-preferred-replica-election.sh" : {},
      "lib/kafka/bin/kafka-producer-perf-test.sh" : {},
      "lib/kafka/bin/kafka-reassign-partitions.sh" : {},
      "lib/kafka/bin/kafka-replay-log-producer.sh" : {},
      "lib/kafka/bin/kafka-replica-verification.sh" : {},
      "lib/kafka/bin/kafka-run-class.sh" : {},
      "lib/kafka/bin/kafka-sentry.sh" : {},
      "lib/kafka/bin/kafka-server-start.sh" : {},
      "lib/kafka/bin/kafka-server-stop.sh" : {},
      "lib/kafka/bin/kafka-simple-consumer-shell.sh" : {},
      "lib/kafka/bin/kafka-streams-application-reset.sh" : {},
      "lib/kafka/bin/kafka-topics.sh" : {},
      "lib/kafka/bin/kafka-verifiable-consumer.sh" : {},
      "lib/kafka/bin/kafka-verifiable-producer.sh" : {},
      "lib/kafka/bin/trogdor.sh" : {},
      "lib/kafka/bin/zookeeper-security-migration.sh" : {},
      "lib/kafka/bin/zookeeper-server-start.sh" : {},
      "lib/kafka/bin/zookeeper-server-stop.sh" : {},
      "lib/kafka/bin/zookeeper-shell.sh" : {},
      "lib/kafka/cloudera" : {},
      "lib/kafka/cloudera/cdh_version.properties" : {},
      "lib/kafka/config" : {},
      "lib/kafka/libs" : {},
      "lib/kafka/libs/activation-1.1.jar" : {},
      "lib/kafka/libs/aopalliance-1.0.jar" : {},
      "lib/kafka/libs/aopalliance-repackaged-2.5.0-b32.jar" : {},
      "lib/kafka/libs/apacheds-i18n-2.0.0-M15.jar" : {},
      "lib/kafka/libs/apacheds-jdbm1-2.0.0-M2.jar" : {},
      "lib/kafka/libs/apacheds-kerberos-codec-2.0.0-M15.jar" : {},
      "lib/kafka/libs/api-asn1-api-1.0.0-M20.jar" : {},
      "lib/kafka/libs/api-util-1.0.0-M20.jar" : {},
      "lib/kafka/libs/argparse4j-0.7.0.jar" : {},
      "lib/kafka/libs/asm-3.1.jar" : {},
      "lib/kafka/libs/avro-1.7.6-cdh5.14.2.jar" : {},
      "lib/kafka/libs/cglib-2.2.1-v20090111.jar" : {},
      "lib/kafka/libs/commons-beanutils-1.8.3.jar" : {},
      "lib/kafka/libs/commons-beanutils-core-1.8.0.jar" : {},
      "lib/kafka/libs/commons-cli-1.2.jar" : {},
      "lib/kafka/libs/commons-codec-1.9.jar" : {},
      "lib/kafka/libs/commons-collections-3.2.2.jar" : {},
      "lib/kafka/libs/commons-compress-1.4.1.jar" : {},
      "lib/kafka/libs/commons-configuration-1.6.jar" : {},
      "lib/kafka/libs/commons-digester-1.8.jar" : {},
      "lib/kafka/libs/commons-el-1.0.jar" : {},
      "lib/kafka/libs/commons-httpclient-3.1.jar" : {},
      "lib/kafka/libs/commons-io-2.4.jar" : {},
      "lib/kafka/libs/commons-lang-2.6.jar" : {},
      "lib/kafka/libs/commons-lang3-3.5.jar" : {},
      "lib/kafka/libs/commons-logging-1.2.jar" : {},
      "lib/kafka/libs/commons-math3-3.1.1.jar" : {},
      "lib/kafka/libs/commons-net-3.1.jar" : {},
      "lib/kafka/libs/commons-pool2-2.4.2.jar" : {},
      "lib/kafka/libs/connect-api-1.0.1-kafka-3.1.0-SNAPSHOT.jar" : {},
      "lib/kafka/libs/connect-file-1.0.1-kafka-3.1.0-SNAPSHOT.jar" : {},
      "lib/kafka/libs/connect-json-1.0.1-kafka-3.1.0-SNAPSHOT.jar" : {},
      "lib/kafka/libs/connect-runtime-1.0.1-kafka-3.1.0-SNAPSHOT.jar" : {},
      "lib/kafka/libs/connect-transforms-1.0.1-kafka-3.1.0-SNAPSHOT.jar" : {},
      "lib/kafka/libs/gson-2.2.4.jar" : {},
      "lib/kafka/libs/guava-20.0.jar" : {},
      "lib/kafka/libs/guice-3.0.jar" : {},
      "lib/kafka/libs/guice-servlet-3.0.jar" : {},
      "lib/kafka/libs/hadoop-annotations-2.6.0-cdh5.14.2.jar" : {},
      "lib/kafka/libs/hadoop-archives-2.6.0-cdh5.14.2.jar" : {},
      "lib/kafka/libs/hadoop-auth-2.6.0-cdh5.14.2.jar" : {},
      "lib/kafka/libs/hadoop-common-2.6.0-cdh5.14.2.jar" : {},
      "lib/kafka/libs/hadoop-mapreduce-client-common-2.6.0-cdh5.14.2.jar" : {},
      "lib/kafka/libs/hadoop-mapreduce-client-core-2.6.0-cdh5.14.2.jar" : {},
      "lib/kafka/libs/hadoop-mapreduce-client-jobclient-2.6.0-cdh5.14.2.jar" : {},
      "lib/kafka/libs/hadoop-mapreduce-client-shuffle-2.6.0-cdh5.14.2.jar" : {},
      "lib/kafka/libs/hadoop-yarn-api-2.6.0-cdh5.14.2.jar" : {},
      "lib/kafka/libs/hadoop-yarn-client-2.6.0-cdh5.14.2.jar" : {},
      "lib/kafka/libs/hadoop-yarn-common-2.6.0-cdh5.14.2.jar" : {},
      "lib/kafka/libs/hadoop-yarn-server-common-2.6.0-cdh5.14.2.jar" : {},
      "lib/kafka/libs/hadoop-yarn-server-nodemanager-2.6.0-cdh5.14.2.jar" : {},
      "lib/kafka/libs/hive-hcatalog-core-1.1.0-cdh5.14.2.jar" : {},
      "lib/kafka/libs/hive-hcatalog-server-extensions-1.1.0-cdh5.14.2.jar" : {},
      "lib/kafka/libs/hk2-api-2.5.0-b32.jar" : {},
      "lib/kafka/libs/hk2-locator-2.5.0-b32.jar" : {},
      "lib/kafka/libs/hk2-utils-2.5.0-b32.jar" : {},
      "lib/kafka/libs/htrace-core4-4.0.1-incubating.jar" : {},
      "lib/kafka/libs/httpclient-4.4.1.jar" : {},
      "lib/kafka/libs/httpcore-4.4.1.jar" : {},
      "lib/kafka/libs/jackson-annotations-2.9.0.jar" : {},
      "lib/kafka/libs/jackson-annotations-2.9.1.jar" : {},
      "lib/kafka/libs/jackson-core-2.9.1.jar" : {},
      "lib/kafka/libs/jackson-core-asl-1.9.2.jar" : {},
      "lib/kafka/libs/jackson-databind-2.9.1.jar" : {},
      "lib/kafka/libs/jackson-jaxrs-1.8.8.jar" : {},
      "lib/kafka/libs/jackson-jaxrs-base-2.9.1.jar" : {},
      "lib/kafka/libs/jackson-jaxrs-json-provider-2.9.1.jar" : {},
      "lib/kafka/libs/jackson-mapper-asl-1.9.2.jar" : {},
      "lib/kafka/libs/jackson-module-jaxb-annotations-2.9.1.jar" : {},
      "lib/kafka/libs/jackson-xc-1.8.8.jar" : {},
      "lib/kafka/libs/jasper-compiler-5.5.23.jar" : {},
      "lib/kafka/libs/jasper-runtime-5.5.23.jar" : {},
      "lib/kafka/libs/java-xmlbuilder-0.4.jar" : {},
      "lib/kafka/libs/javassist-3.20.0-GA.jar" : {},
      "lib/kafka/libs/javassist-3.21.0-GA.jar" : {},
      "lib/kafka/libs/javax.annotation-api-1.2.jar" : {},
      "lib/kafka/libs/javax.inject-1.jar" : {},
      "lib/kafka/libs/javax.inject-2.5.0-b32.jar" : {},
      "lib/kafka/libs/javax.servlet-api-3.1.0.jar" : {},
      "lib/kafka/libs/javax.ws.rs-api-2.0.1.jar" : {},
      "lib/kafka/libs/jaxb-api-2.2.2.jar" : {},
      "lib/kafka/libs/jersey-client-2.25.1.jar" : {},
      "lib/kafka/libs/jersey-common-2.25.1.jar" : {},
      "lib/kafka/libs/jersey-container-servlet-2.25.1.jar" : {},
      "lib/kafka/libs/jersey-container-servlet-core-2.25.1.jar" : {},
      "lib/kafka/libs/jersey-guava-2.25.1.jar" : {},
      "lib/kafka/libs/jersey-guice-1.9.jar" : {},
      "lib/kafka/libs/jersey-media-jaxb-2.25.1.jar" : {},
      "lib/kafka/libs/jersey-server-2.25.1.jar" : {},
      "lib/kafka/libs/jets3t-0.9.0.jar" : {},
      "lib/kafka/libs/jettison-1.1.jar" : {},
      "lib/kafka/libs/jetty-6.1.26.cloudera.4.jar" : {},
      "lib/kafka/libs/jetty-continuation-9.2.22.v20170606.jar" : {},
      "lib/kafka/libs/jetty-http-9.2.22.v20170606.jar" : {},
      "lib/kafka/libs/jetty-io-9.2.22.v20170606.jar" : {},
      "lib/kafka/libs/jetty-security-9.2.22.v20170606.jar" : {},
      "lib/kafka/libs/jetty-server-9.2.22.v20170606.jar" : {},
      "lib/kafka/libs/jetty-servlet-9.2.22.v20170606.jar" : {},
      "lib/kafka/libs/jetty-servlets-9.2.22.v20170606.jar" : {},
      "lib/kafka/libs/jetty-util-6.1.26.cloudera.4.jar" : {},
      "lib/kafka/libs/jetty-util-9.2.22.v20170606.jar" : {},
      "lib/kafka/libs/jopt-simple-5.0.4.jar" : {},
      "lib/kafka/libs/jsch-0.1.42.jar" : {},
      "lib/kafka/libs/jsp-api-2.1.jar" : {},
      "lib/kafka/libs/jsr305-3.0.0.jar" : {},
      "lib/kafka/libs/kafka-clients-1.0.1-kafka-3.1.0-SNAPSHOT.jar" : {},
      "lib/kafka/libs/kafka-log4j-appender-1.0.1-kafka-3.1.0-SNAPSHOT.jar" : {},
      "lib/kafka/libs/kafka-streams-1.0.1-kafka-3.1.0-SNAPSHOT.jar" : {},
      "lib/kafka/libs/kafka-streams-examples-1.0.1-kafka-3.1.0-SNAPSHOT.jar" : {},
      "lib/kafka/libs/kafka-tools-1.0.1-kafka-3.1.0-SNAPSHOT.jar" : {},
      "lib/kafka/libs/kafka_2.11-1.0.1-kafka-3.1.0-SNAPSHOT-sources.jar" : {},
      "lib/kafka/libs/kafka_2.11-1.0.1-kafka-3.1.0-SNAPSHOT-test-sources.jar" : {},
      "lib/kafka/libs/kafka_2.11-1.0.1-kafka-3.1.0-SNAPSHOT.jar" : {},
      "lib/kafka/libs/leveldbjni-all-1.8.jar" : {},
      "lib/kafka/libs/libthrift-0.9.3.jar" : {},
      "lib/kafka/libs/log4j-1.2.17.jar" : {},
      "lib/kafka/libs/lz4-java-1.4.jar" : {},
      "lib/kafka/libs/maven-artifact-3.5.0.jar" : {},
      "lib/kafka/libs/metrics-core-2.2.0.jar" : {},
      "lib/kafka/libs/metrics-servlet-2.2.0.jar" : {},
      "lib/kafka/libs/osgi-resource-locator-1.0.1.jar" : {},
      "lib/kafka/libs/paranamer-2.3.jar" : {},
      "lib/kafka/libs/plexus-utils-3.0.24.jar" : {},
      "lib/kafka/libs/protobuf-java-2.5.0.jar" : {},
      "lib/kafka/libs/reflections-0.9.11.jar" : {},
      "lib/kafka/libs/rocksdbjni-5.7.3.jar" : {},
      "lib/kafka/libs/scala-library-2.11.12.jar" : {},
      "lib/kafka/libs/sentry-binding-hive-conf-1.5.1-cdh5.14.2.jar" : {},
      "lib/kafka/libs/sentry-binding-hive-follower-1.5.1-cdh5.14.2.jar" : {},
      "lib/kafka/libs/sentry-binding-kafka-1.5.1-cdh5.14.2.jar" : {},
      "lib/kafka/libs/sentry-core-common-1.5.1-cdh5.14.2.jar" : {},
      "lib/kafka/libs/sentry-core-model-db-1.5.1-cdh5.14.2.jar" : {},
      "lib/kafka/libs/sentry-core-model-indexer-1.5.1-cdh5.14.2.jar" : {},
      "lib/kafka/libs/sentry-core-model-kafka-1.5.1-cdh5.14.2.jar" : {},
      "lib/kafka/libs/sentry-hdfs-common-1.5.1-cdh5.14.2.jar" : {},
      "lib/kafka/libs/sentry-policy-common-1.5.1-cdh5.14.2.jar" : {},
      "lib/kafka/libs/sentry-policy-indexer-1.5.1-cdh5.14.2.jar" : {},
      "lib/kafka/libs/sentry-policy-kafka-1.5.1-cdh5.14.2.jar" : {},
      "lib/kafka/libs/sentry-provider-common-1.5.1-cdh5.14.2.jar" : {},
      "lib/kafka/libs/sentry-provider-db-1.5.1-cdh5.14.2.jar" : {},
      "lib/kafka/libs/sentry-provider-file-1.5.1-cdh5.14.2.jar" : {},
      "lib/kafka/libs/servlet-api-2.5.jar" : {},
      "lib/kafka/libs/shiro-core-1.2.3.jar" : {},
      "lib/kafka/libs/slf4j-api-1.7.12.jar" : {},
      "lib/kafka/libs/slf4j-api-1.7.5.jar" : {},
      "lib/kafka/libs/slf4j-log4j12-1.7.5.jar" : {},
      "lib/kafka/libs/snappy-java-1.1.4.jar" : {},
      "lib/kafka/libs/stax-api-1.0-2.jar" : {},
      "lib/kafka/libs/validation-api-1.1.0.Final.jar" : {},
      "lib/kafka/libs/xmlenc-0.52.jar" : {},
      "lib/kafka/libs/xz-1.0.jar" : {},
      "lib/kafka/libs/zkclient-0.10.jar" : {},
      "lib/kafka/libs/zookeeper-3.4.5-cdh5.14.2.jar" : {},
      "lib/kafka/site-docs" : {},
      "lib/kafka/site-docs/kafka_2.11-1.0.1-kafka-3.1.0-SNAPSHOT-site-docs.tgz" : {}
    }
  }
}
